# tasty.ini
[program:tasty]
command=/usr/local/bin/usenv /opt/tasty-server/current_env python /opt/tasty-server/current_bundle/src/local.py --port=8080 --config=/etc/tasty.conf
redirect_stderr=true
stdout_logfile=/var/log/tasty-server.log
stdout_logfile_maxbytes=104857600
autorestart=true
startsecs=5
