# tasty.ini
[program:tasty]
command=/usr/local/bin/usenv /opt/echonest/tastyenv python /opt/echonest/tasty-server/src/local.py --port=8080 --config=/opt/echonest/tasty-server/src/tasty.conf
directory=/opt/echonest/tasty-server
redirect_stderr=true
stdout_logfile=/var/log/tasty-server.log
stdout_logfile_maxbytes=104857600
autorestart=true
startsecs=5
